# systemd/digital-secretariat.service.in
[Unit]
Description=Digital Secretariat private control plane
After=network-online.target
Wants=network-online.target
ConditionPathExists=@INSTALL_ROOT@/current/src/secretariat/app.py
ConditionPathExists=@CONFIG_ROOT@/secretariat.env
ConditionPathExists=@CONFIG_ROOT@/database-url
ConditionPathExists=@CONFIG_ROOT@/action-database-url
ConditionPathExists=@CONFIG_ROOT@/action-token-secret

[Service]
Type=simple
User=@APP_USER@
Group=@APP_GROUP@
WorkingDirectory=@INSTALL_ROOT@/current
Environment=PYTHONPATH=@INSTALL_ROOT@/current
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=SECRETARIAT_DATABASE_URL_FILE=/run/credentials/digital-secretariat.service/database-url
Environment=SECRETARIAT_SNAPSHOT_FILE=/var/cache/digital-secretariat/dashboard.json
Environment=SECRETARIAT_MAIL_CACHE=/var/cache/digital-secretariat-mail
Environment=MAIL_RAIL_DATABASE_URL_FILE=/run/credentials/digital-secretariat.service/action-database-url
Environment=DATABASE_URL_FILE=/run/credentials/digital-secretariat.service/action-database-url
Environment=SECRETARIAT_ACTION_SECRET_FILE=/run/credentials/digital-secretariat.service/action-token-secret
EnvironmentFile=@CONFIG_ROOT@/secretariat.env
LoadCredential=database-url:@CONFIG_ROOT@/database-url
LoadCredential=action-database-url:@CONFIG_ROOT@/action-database-url
LoadCredential=action-token-secret:@CONFIG_ROOT@/action-token-secret
CacheDirectory=digital-secretariat
CacheDirectoryMode=0700
ExecStart=@PYTHON@ -m uvicorn src.secretariat.app:app --host 127.0.0.1 --port @PORT@ --no-access-log --no-proxy-headers
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=read-only
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
SystemCallArchitectures=native
CapabilityBoundingSet=
AmbientCapabilities=
IPAddressDeny=any
IPAddressAllow=localhost
IPAddressAllow=100.100.100.100
IPAddressAllow=fd7a:115c:a1e0::53
@SMTP_IP_ALLOW_RULES@
UMask=0077

StandardOutput=journal
StandardError=journal
SyslogIdentifier=digital-secretariat

[Install]
WantedBy=multi-user.target
