# systemd/digital-secretariat-mail-cache.service.in
[Unit]
Description=Digital Secretariat read-only mail archive refresh
After=network-online.target
Wants=network-online.target
ConditionPathExists=@INSTALL_ROOT@/current/src/secretariat/mail_content.py
ConditionPathExists=@CONFIG_ROOT@/content-database-url

[Service]
Type=oneshot
User=@APP_USER@
Group=@APP_GROUP@
WorkingDirectory=@INSTALL_ROOT@/current
Environment=PYTHONPATH=@INSTALL_ROOT@/current
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=SECRETARIAT_CONTENT_DATABASE_URL_FILE=/run/credentials/digital-secretariat-mail-cache.service/database-url
Environment=SECRETARIAT_MAIL_CACHE=/var/cache/digital-secretariat-mail
LoadCredential=database-url:@CONFIG_ROOT@/content-database-url
CacheDirectory=digital-secretariat-mail
CacheDirectoryMode=0700
ExecStart=@PYTHON@ -m src.secretariat.mail_content --limit 200
TimeoutStartSec=15min

NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=read-only
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
SystemCallArchitectures=native
CapabilityBoundingSet=
AmbientCapabilities=
UMask=0077

StandardOutput=journal
StandardError=journal
SyslogIdentifier=digital-secretariat-mail-cache
