# systemd/digital-secretariat-mail-cache.timer
[Unit]
Description=Refresh private Secretariat mail archive every five minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=true
Unit=digital-secretariat-mail-cache.service

[Install]
WantedBy=timers.target
