# systemd/tailscaled-secretariat.service
[Unit]
Description=Isolated Tailscale identity for the Digital Secretariat
Documentation=https://tailscale.com/docs/features/tailscale-services
Wants=network-online.target
After=network-online.target

[Service]
Type=notify
NotifyAccess=main
DynamicUser=true
RuntimeDirectory=tailscale-secretariat
RuntimeDirectoryMode=0700
StateDirectory=tailscale-secretariat
StateDirectoryMode=0700
CacheDirectory=tailscale-secretariat
CacheDirectoryMode=0700
ExecStart=/usr/sbin/tailscaled --state=/var/lib/tailscale-secretariat/tailscaled.state --statedir=/var/lib/tailscale-secretariat --socket=/run/tailscale-secretariat/tailscaled.sock --tun=userspace-networking --port=0
Restart=on-failure
RestartSec=5

NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectHostname=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
CapabilityBoundingSet=
AmbientCapabilities=
UMask=0077

StandardOutput=journal
StandardError=journal
SyslogIdentifier=tailscaled-secretariat

[Install]
WantedBy=multi-user.target
